# fcuusbhost.service
[Unit]
Description=FCU Usb Service
 
[Service]
Type=notify
ExecStart=/home/demis/FCUUsbHost
 
[Install]
WantedBy=multi-user.target
